# components/cookbooks/keywhiz-client/templates/default/secrets.mount.erb
# Keysync secrets tmpfs mount unit file
#
# Auto generated by OneOps
# Created On : <%= Time.now.getutc %>

[Unit]
Description=Keysync secrets tmpfs mount point.

[Mount]
What=tmpfs
Where=<%= @secrets_dir %>
Type=tmpfs
Options=uid=root,gid=root,nodev,nosuid,noexec,size=100M,nr_inodes=100k,mode=0755

[Install]
WantedBy=multi-user.target